# build/live-build/config/includes.chroot/etc/xdg/autostart/senior-zero-branding.desktop
[Desktop Entry]
Type=Application
Name=Aegis Session Branding
Comment=Apply Aegis visual branding at session start
Exec=senior-zero-branding-apply
OnlyShowIn=LXDE;
X-GNOME-Autostart-enabled=true
